# configs/nyc/files/etc/uci-defaults/01_system
#!/bin/sh


uci batch <<EOF
set system.@system[0].hostname='chelsea'
set system.@system[0].timezone='EST5EDT,M3.2.0,M11.1.0'
set system.@system[0].zonename='America/New York'
set system.@system[0].log_size='1024'
EOF

uci batch <<EOF
set system.status='led'
set system.status.name='power'
set system.status.sysfs='green:power'
set system.status.trigger='none'
set system.status.default='1'
EOF

uci delete system.led_lan
uci batch <<EOF
set system.wan='led'
set system.wan.name='WAN'
set system.wan.sysfs='green:mesh'
set system.wan.trigger='netdev'
set system.wan.dev='eth1'
set system.wan.mode='link tx rx'
EOF

uci delete system.led_wlan
uci batch <<EOF
set system.wlan='led'
set system.wlan.name='WLAN'
set system.wlan.sysfs='green:wlan'
set system.wlan.trigger='phy1tpt'
EOF

uci commit system

exit 0
