# configs/mikrotik/files/etc/uci-defaults/05_wan
#!/bin/sh

uci batch <<EOF
set network.wan.disabled='1'
set network.wan6.disabled='1'
EOF
uci commit network

exit 0
