# configs/dacha/files/etc/uci-defaults/05_wan
#!/bin/sh


uci batch <<EOF
set network.wan.ifname='usb0'
set network.wan.proto='dhcp'
EOF
uci commit network

exit 0
