# configs/dacha/files/etc/uci-defaults/06_wan
#!/bin/sh

uci batch <<EOF
set network.wan='interface'
set network.wan.ifname='eth1'
set network.wan.proto='dhcp'
EOF
uci commit network

exit 0
