# configs/srkrt/files/etc/uci-defaults/01_system
#!/bin/sh


uci batch <<EOF
set system.@system[0].hostname='srkrt'
set system.@system[0].timezone='UTC'
set system.@system[0].log_size='1024'
EOF
uci commit system

exit 0
